# tags/7.0.3/debian/postinst
#!/bin/sh -e

if [ "$1" = "configure" ]; then
	if [ -z "$2" ]; then
		update-rc.d oops-firewall defaults 99 10 > /dev/null
	fi
fi
