# tags/6.2.7/debian/prerm
#!/bin/sh

set -e
if [ "$1" = "purge" ]; then
	invoke-rc.d oops-firewall stop --quite > /dev/null
fi
